# etc/cwput.conf
description "put cloudwatch metrics"

start on runlevel [2]
stop on runlevel [016]

script
    /usr/bin/cwput.start
end script
